# DockerLogDriver/src/main/resources/dev.socket
[Unit]
Description=Your plugin

[Socket]
ListenStream=/run/docker/plugins/dev.sock

[Install]
WantedBy=sockets.target